# systemd/allchat.service
[Unit]
Description=All_Chat - FastAPI Application Server
Documentation=https://github.com/YOUR_USERNAME/all_chat
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=allchat
Group=allchat
WorkingDirectory=/app/backend

EnvironmentFile=/app/backend/.env
Environment="PATH=/app/venv/bin:/usr/local/bin:/usr/bin:/bin"

# 1 worker — safe default for modest servers.
# Increase to 2 or 4 on boxes with 2GB+ RAM.
ExecStart=/app/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --loop uvloop \
    --http httptools \
    --no-access-log \
    --proxy-headers \
    --forwarded-allow-ips="127.0.0.1"

ExecReload=/bin/kill -s HUP $MAINPID

Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/app/media /app/backend /tmp

# Resource limits — appropriate for 1-2GB RAM servers
LimitNOFILE=65536
LimitNPROC=512
TasksMax=128
MemoryMax=512M
CPUQuota=90%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=allchat

[Install]
WantedBy=multi-user.target
