# systemd/allchat.service
[Unit]
Description=All_Chat - FastAPI Application Server
Documentation=https://github.com/YOUR_USERNAME/all_chat
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=allchat
Group=allchat
WorkingDirectory=/app/backend

# Environment
EnvironmentFile=/app/backend/.env
Environment="PATH=/app/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Start command — uvicorn with multiple workers
ExecStart=/app/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 4 \
    --loop uvloop \
    --http httptools \
    --no-access-log \
    --proxy-headers \
    --forwarded-allow-ips="127.0.0.1"

# Graceful reload (sends SIGHUP)
ExecReload=/bin/kill -s HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/app/media /app/backend /tmp
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
TasksMax=256
MemoryMax=1G
CPUQuota=200%

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=allchat

[Install]
WantedBy=multi-user.target
